# repos/zarazhangrui--youtube-to-ebook/runs/2026-05-05/run-static-checks/artifacts/com.youtube.newsletter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.youtube.newsletter</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/bytedance/youtube-newsletter/run_newsletter.sh</string>
    </array>

    <!-- Run every Wednesday at 7:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>3</integer>  <!-- 3 = Wednesday -->
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/bytedance/youtube-newsletter/logs/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/bytedance/youtube-newsletter/logs/launchd_error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/bytedance/youtube-newsletter</string>
</dict>
</plist>
